# services/slow-log-cron/config/cron/crontab
* * * * * bash /commands/query-digest.sh
