# services/slow-log-cron/config/cron/crontab
5 * * * * bash /commands/query-digest.sh
